# systemd/reconraven-api.service
[Unit]
Description=ReconRaven REST API Server
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/ReconRaven
ExecStart=/usr/bin/python3 /home/pi/ReconRaven/api/server.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target

